# files/munin/dnf.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=1d
